# deploy/attentionx-metadata.service
[Unit]
Description=AttentionX RISE Metadata Server (port 3006)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=attentionx
Group=attentionx
WorkingDirectory=/opt/attentionx/backend
ExecStart=/usr/bin/node server.js
Environment=NODE_ENV=production

Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

KillSignal=SIGTERM
TimeoutStopSec=10
NoNewPrivileges=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=attentionx-metadata

[Install]
WantedBy=multi-user.target
